# .config/systemd/user/picom.service
# path:   /home/klassiker/.local/share/repos/dotfiles/.config/systemd/user/picom.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/dotfiles
# date:   2021-01-15T13:07:12+0100

[Unit]
Description=picom compositor

[Service]
Type=forking
ExecStart=/usr/bin/picom -b

[Install]
WantedBy=i3.target
