# .config/systemd/user/picom.service
# path:       /home/klassiker/.config/systemd/user/picom.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-10-29T14:54:37+0100

[Unit]
Description=picom compositor
After=screenlayout.service

[Service]
Type=forking
ExecStart=/usr/bin/picom -b

[Install]
WantedBy=i3.target
