# .config/systemd/user/polybar.service
[Unit]
Description=polybar status bars
PartOf=default.target
After=Screenlayout.service

[Service]
Type=forking
ExecStart=%h/repos/polybar/polybar.sh

[Install]
WantedBy=window-manager.target
