# .config/systemd/user/polybar.service
# path:       /home/klassiker/.config/systemd/user/polybar.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-10-29T18:12:37+0100

[Unit]
Description=polybar status bars
After=screenlayout.service

[Service]
Type=forking
ExecStart=%h/.local/share/repos/polybar/polybar.sh

[Install]
WantedBy=i3.target
