# .config/systemd/user/polybar.service
# path:       /home/klassiker/.config/systemd/user/polybar.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-04-30T18:28:16+0200

[Unit]
Description=polybar status bars
PartOf=default.target
After=Screenlayout.service

[Service]
Type=forking
ExecStart=%h/.local/share/repos/polybar/polybar.sh

[Install]
WantedBy=window-manager.target
