# .config/systemd/user/rss.timer
[Unit]
Description=run rss.service every 30 minutes

[Timer]
OnActiveSec=30s
OnUnitActiveSec=30m
AccuracySec=5m
Unit=rss.service

[Install]
WantedBy=timers.target
