# .config/systemd/user/screenlayout.service
[Unit]
Description=setup screen layouts
PartOf=default.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/repos/system/screenlayout.sh

[Install]
WantedBy=window-manager.target
