# .config/systemd/user/window-manager.target
[Unit]
Description=window manager session
BindsTo=default.target
After=window-manager-pre.target
