# .config/systemd/user/xautolock.service
# path:       /home/klassiker/.config/systemd/user/xautolock.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-09-19T09:00:22+0200

[Unit]
Description=xautolock to auto lock/suspend device
PartOf=default.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/xautolock -detectsleep
Restart=on-failure
RestartSec=1m

[Install]
WantedBy=i3.target
