# .config/systemd/user/xautolock.service
# path:   /home/klassiker/.config/systemd/user/xautolock.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/dotfiles
# date:   2021-05-28T10:09:15+0200

[Unit]
Description=xautolock to auto lock/suspend device

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/xautolock -detectsleep
Restart=on-failure
RestartSec=1m

[Install]
WantedBy=i3.target
