# .config/systemd/user/xbanish.service
# path:   /home/klassiker/.config/systemd/user/xbanish.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/dotfiles
# date:   2021-05-28T10:09:19+0200

[Unit]
Description=xbanish to remove mouse cursor on keypress

[Service]
Type=simple
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=i3.target
