# .config/systemd/user/xbanish.service
# path:       /home/klassiker/.config/systemd/user/xbanish.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-10-29T14:53:44+0100

[Unit]
Description=xbanish to remove mouse cursor on keypress

[Service]
Type=simple
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=i3.target
