# .config/systemd/user/xwallpaper.service
# path:       /home/klassiker/.config/systemd/user/xwallpaper.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/dotfiles
# date:       2020-11-18T17:16:16+0100

[Unit]
Description=set wallpaper with xwallpaper
After=screenlayout.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/share/repos/shell/xwallpaper.sh

[Install]
WantedBy=i3.target
