# .config/systemd/user/xwallpaper.service
[Unit]
Description=set wallpaper with xwallpaper
PartOf=default.target
After=screenlayout.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/xwallpaper --zoom %h/.local/share/wallpaper.jpg

[Install]
WantedBy=window-manager.target
